# etc/service/junto-api/run
#!/bin/sh
PORT=4000 exec /opt/app/api/bin/api start
