# etc/service/junto-api/run
#!/bin/sh
exec /opt/app/api/bin/api start
