# etc/service/junto-frontend/run
#!/bin/sh
export PORT=3000
cd /opt/app/frontend
exec npm run start
