# etc/service/junto-frontend/run
#!/bin/sh
cd /opt/app/frontend
exec npm run start
